# scripts/pegasus.service
[Unit]
Description=Pegasus Autonomous AI Agent
Documentation=https://github.com/faelsete/pegasus-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pegasus-agent
ExecStart=/usr/bin/node --import tsx/esm src/index.ts
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
Environment=NODE_ENV=production
Environment=HOME=/root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pegasus

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=true

# Resource limits
TasksMax=256

[Install]
WantedBy=multi-user.target
